# otmonitor/rootfs/etc/services.d/otmonitor/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Bas Nijholt's Hass.io Add-ons: otmonitor
# Runs the otmonitor server
# ==============================================================================

bashio::log.info "Starting otmonitor..."

# RUN otmonitor
cd /opt || exit
sleep 86400
exec ./otmonitor --daemon -f /etc/otmonitor/otmonitor.conf
